# systemd/cascade-walker@.service
[Unit]
Description=Cascade Walker — Phase 2 sub-engagement walker for tweet %i
After=network.target tweet-tracker@%i.service
Wants=tweet-tracker@%i.service

[Service]
Type=simple
WorkingDirectory=/opt/tweet-tracker
Environment=TWITTER241_RAPIDAPI_KEY=YOUR_KEY_HERE
Environment=TWEET_ID=%i
Environment=DATA_DIR=/opt/tweet-tracker/data
Environment=WALKER_INTERVAL_SEC=1800
ExecStart=/usr/bin/python3 /opt/tweet-tracker/cascade_walker.py
Restart=always
RestartSec=30
StandardOutput=append:/opt/tweet-tracker/walker.log
StandardError=append:/opt/tweet-tracker/walker.log

[Install]
WantedBy=multi-user.target
