# systemd/cascade-walker@.service
[Unit]
Description=Cascade Walker — Phase 2 sub-engagement walker for tweet %i
After=network.target tweet-tracker@%i.service
Wants=tweet-tracker@%i.service

[Service]
Type=simple
WorkingDirectory=/opt/tweet-tracker
Environment=XAPI_API_KEY=sk-ae375f79914f27c6e29ec9aaad09448378b133df175bcac4
Environment=TWEET_ID=%i
Environment=DATA_DIR=/opt/tweet-tracker/data
Environment=WALKER_INTERVAL_SEC=1800
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /opt/tweet-tracker/cascade_walker.py
Restart=always
RestartSec=30
StandardOutput=append:/opt/tweet-tracker/walker.log
StandardError=append:/opt/tweet-tracker/walker.log

[Install]
WantedBy=multi-user.target
