# systemd/tweet-tracker@.service
[Unit]
Description=XHI Tweet Tracker — high-frequency snapshot of tweet %i
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/tweet-tracker
Environment=TWITTER241_RAPIDAPI_KEY=YOUR_KEY_HERE
Environment=TWEET_ID=%i
Environment=DATA_DIR=/opt/tweet-tracker/data
Environment=SNAPSHOT_INTERVAL_SEC=300
Environment=MAX_PAGES_PER_CYCLE=5
ExecStart=/usr/bin/python3 /opt/tweet-tracker/tracker.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/tweet-tracker/tracker.log
StandardError=append:/opt/tweet-tracker/tracker.log

[Install]
WantedBy=multi-user.target
