# systemd/tweet-tracker@.service
[Unit]
Description=XHI Tweet Tracker — high-frequency snapshot of tweet %i
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/tweet-tracker
Environment=XAPI_API_KEY=sk-ae375f79914f27c6e29ec9aaad09448378b133df175bcac4
Environment=TWEET_ID=%i
Environment=DATA_DIR=/opt/tweet-tracker/data
Environment=SNAPSHOT_INTERVAL_SEC=300
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /opt/tweet-tracker/tracker.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/tweet-tracker/tracker.log
StandardError=append:/opt/tweet-tracker/tracker.log

[Install]
WantedBy=multi-user.target
